# evil-sensor/evil-sensor@.service
[Unit]
Description=evil sensor on %I interface

[Service]
ExecStart=-/home/pi/evil-sensor.py %I
Restart=always
RestartSec=5
